# infra/vast/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[program:vllm]
command=/bin/bash /app/infra/vast/run-vllm.sh
directory=/app
autostart=true
autorestart=unexpected
startsecs=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=10

[program:redis]
command=/bin/bash /app/infra/vast/run-redis.sh
directory=/app
autostart=true
autorestart=unexpected
startsecs=2
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=20

[program:api]
command=/bin/bash /app/infra/vast/run-api.sh
directory=/app
autostart=true
autorestart=true
startsecs=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=30

[program:cloudflared]
command=/bin/bash /app/infra/vast/run-cloudflared.sh
directory=/app
autostart=true
autorestart=unexpected
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
priority=40
