# infra/supervisord.conf
[supervisord]
nodaemon=true

[program:api]
command=uv run uvicorn backend.main:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true

